# conf/openresty.conf
[program:openresty]
command=/usr/local/bin/nginx -c /etc/nginx/nginx.conf
autorestart=true